# playbooks/files/deb_perso/user-mount_systemd/etc/user-mount/template_user.automount
[Unit]
Description=Partages ${NAME_SHARE} ${USER_LOCAL}
Requires=network-online.target

[Automount]
Where=${PATH_MOUNT}
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target

